# deploy/tg-music-soundcloud.service
[Unit]
Description=Telegram Music Bot — SoundCloud import worker (Celery, queue=soundcloud)
After=network-online.target redis-server.service
Wants=network-online.target
Requires=redis-server.service

[Service]
WorkingDirectory=/opt/tg-music-bot
# Своя очередь и свой воркер — не делит concurrency с youtube (там бэклог из
# сотен поштучных импортов). concurrency=3 — темп 3000 треков/сутки
# (SPEC-КАТАЛОГ §4): 7 прокси в ротации размазывают параллельные запросы по
# разным IP, анти-бан-паузы остаются на каждой задаче. Без прокси вернуть 1.
ExecStart=/opt/tg-music-bot/.venv/bin/celery -A app.tasks.celery_app worker -Q soundcloud --pool=threads --concurrency=3 --loglevel=info
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
