# deploy/tg-music-telegram-channel-scan.timer
[Unit]
Description=Daily check for the Telegram channel due for re-scan

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
